# .claude/runbooks/ai.hermes.monitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.hermes.monitor</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/yuseungju/AIJOB/scripts/hermes_monitor.py</string>
      <string>--write-snapshot</string>
      <string>/Users/yuseungju/AIJOB/.claude/tmp/hermes-monitor/latest.json</string>
      <string>--notify</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/yuseungju/AIJOB</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>HERMES_MONITOR_SLACK_WEBHOOK_URL</key>
      <string>REPLACE_ME</string>
    </dict>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/yuseungju/Library/Logs/hermes-monitor.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/yuseungju/Library/Logs/hermes-monitor.error.log</string>
  </dict>
</plist>
